# deploy/systemd/dierchat.service.example
# /etc/systemd/system/dierchat.service
# Пути и пользователя подставьте под свой сервер.
#
# sudo systemctl daemon-reload
# sudo systemctl enable --now dierchat

[Unit]
Description=DierCHAT API + web (Go)
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/dierchat
ExecStart=/opt/dierchat/dierchat-server
Restart=on-failure
Environment=DIERCHAT_CONFIG=/opt/dierchat/config.json

[Install]
WantedBy=multi-user.target
